# rootfs/etc/services.d/storjshare-daemon/run
#!/usr/bin/execlineb -P

/usr/bin/storjshare-daemon -F
