# rootfs/etc/services.d/02-storjshare-start/run
#!/usr/bin/execlineb -P

storjshare-start -c "/config/config.json"
